# crontab
* * * * * root echo "job done" >> /log.txt
